# ubuntu/14.04/message-consumer/mc.supervisord.conf
[program:mc]
command=/docker-entrypoint.sh
stdout_logfile=/var/log/message-consumer/message-consumer.log
stderr_logfile=/var/log/message-consumer/message-consumer.log
