# ubuntu/14.04/message-consumer/message-consumer.supervisord.conf
[program:message-consumer]
command=/run-message-consumer.sh
stdout_logfile=/var/log/message-consumer/message-consumer.log
stderr_logfile=/var/log/message-consumer/message-consumer.log
