# ubuntu/14.04/gluuoxeleven/jetty.conf
[program:jetty]
# switched to jetty base directory
directory=/opt/jetty_base/oxeleven
command=java -jar /opt/jetty/start.jar
