# examples/multi-hosts/csync2-manager.cron
1-59/2 * * * *    root    /usr/sbin/csync2 -N manager-1.gluu -xv 2>/var/log/csync2.log
