# examples/multi-hosts/csync2.cron.tmpl
2-59/2 * * * *    root    /usr/sbin/csync2 -N NODE.gluu -xv 2>/var/log/csync2.log
