# ubuntu/14.04/gluu-flask-dockerized/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=2
loglevel=error ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid
user=root
childlogdir=/var/log/supervisor/

[program:run]
command=bash -c ". /root/.bashrc; /root/gluu-flask/run.py"

[program:master]
command=/usr/bin/salt-master
