# ubuntu/14.04/oxeleven/tomcat.supervisord.conf
[program:tomcat]
command=/opt/tomcat/bin/catalina.sh run
environment=CATALINA_PID=/var/run/tomcat.pid
