# ubuntu/14.04/_oxasimba/jetty.supervisord.conf
[program:jetty]
directory=/opt/gluu/jetty/asimba
command=java -jar /opt/jetty/start.jar -server -Xms256m -Xmx1024m -XX:+DisableExplicitGC -Dgluu.base=/etc/gluu -Dcatalina.base=/opt/gluu/jetty/asimba
